# systemd/roberthalf-scraper.service
[Unit]
Description=Robert Half Job Scraper (Run Once, Triggered by Timer)
# This service is triggered by roberthalf-scraper.timer
After=network.target

[Service]
Type=simple
User=jstocks
Group=jstocks
WorkingDirectory=/home/jstocks/PROJECTS/scrape-roberthalf

# Environment setup
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/jstocks/PROJECTS/scrape-roberthalf
EnvironmentFile=/home/jstocks/PROJECTS/scrape-roberthalf/.env

# Execution
ExecStart=/home/jstocks/.local/bin/uv run python roberthalf_scraper.py

# Restart configuration
Restart=on-failure
RestartSec=30
StartLimitInterval=1800
StartLimitBurst=5

# Resource limits
CPUQuota=50%
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
